# debian/prey.cron.d
# /etc/cron.d/prey: crontab entries for the prey package

*/20 * * * * root $INSTALLPATH/$PREY_FILE >/dev/null
